# supervisord/configs/payload.conf
[supervisord]

[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /home/ec2-user/projects/derf/uwsgi/payload.ini
user=ec2-user

